# supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:fastapi]
command=uvicorn app.main:app --host 127.0.0.1 --port 8001
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fastapi-out.log
stderr_logfile=/var/log/supervisor/fastapi-err.log
environment=PANSEARCH_DATA_DIR="/app/data"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-out.log
stderr_logfile=/var/log/supervisor/nginx-err.log
